# conf/translate-service.service
[Unit]
Description=My FastAPI app
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/path/to/myapp
ExecStart=/usr/bin/env python3 -m uvicorn main:app --host 0.0.0.0 --port 8080
Restart=always

[Install]
WantedBy=multi-user.target